# agents-sdk/schedules/com.sean.agent.gemini-researcher.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sean.agent.gemini-researcher</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/seanwinslow/Code-Brain/code-brain/agents-sdk/.venv/bin/python3</string>
        <string>/Users/seanwinslow/Code-Brain/code-brain/agents-sdk/agents/gemini_researcher.py</string>
        <string>--mode</string>
        <string>queue</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/seanwinslow/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/seanwinslow/Code-Brain/code-brain</string>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/seanwinslow/Code-Brain/code-brain/vault/90_system/agent-logs/gemini-researcher-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seanwinslow/Code-Brain/code-brain/vault/90_system/agent-logs/gemini-researcher-stderr.log</string>
</dict>
</plist>
