# agents-sdk/schedules/com.sean.agent.daily-morning-baton.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
  Baton-triggered Daily Driver: fires when Process Inbox completes.
  WatchPaths monitors the inbox_done.flag baton file.
  This supplements (not replaces) the time-based daily-morning plist.
  If the baton already triggered the morning run, the time-based run
  will find the daily note already populated and add minimal content.
-->
<dict>
    <key>Label</key>
    <string>com.sean.agent.daily-morning-baton</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/seanwinslow/Code-Brain/claude-code-superuser-pack/agents-sdk/.venv/bin/python3</string>
        <string>/Users/seanwinslow/Code-Brain/claude-code-superuser-pack/agents-sdk/agents/daily_driver.py</string>
        <string>--mode</string>
        <string>morning</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/Users/seanwinslow/.claude/batons/inbox_done.flag</string>
    </array>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/seanwinslow/.local/bin:/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/seanwinslow/Code-Brain/claude-code-superuser-pack</string>
    <key>StandardOutPath</key>
    <string>/Users/seanwinslow/Code-Brain/claude-code-superuser-pack/vault/90_system/agent-logs/daily-morning-baton-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/seanwinslow/Code-Brain/claude-code-superuser-pack/vault/90_system/agent-logs/daily-morning-baton-stderr.log</string>
</dict>
</plist>
